# brew-installer/brew-service.server.service
[Unit]
Description=Brew Startup

[Service]
Type=simple
ExecStart=bash -c "sudo -H -u a bash brew.service.sh"
 
[Install]
WantedBy=multi-user.target