# brew-installer/brew-service.server.service
[Unit]
Description=Brew Startup

[Service]
Type=any
ExecStart=bash /home/a/git-repos/common-linux-installer/brew-installer/brew.service.sh

[Install]
WantedBy=multi-user.target